# custom-services/autostart.service
[Install]
WantedBy=autostart.target

[Unit]
Description=<service description>

[Service]
ExecStart=/usr/bin/dex -a ~/.config/autostart/
