# custom-services/kanata.service
[Unit]
Description=kanata

[Service]
Type=simple
ExecStart=kanata -c /etc/kanata/kanata.kbd -p 5829 -n
Restart=no

[Install]
WantedBy=sysinit.target
