# custom-services/updategeo.service
[Unit]
Description=Update geosite geoip

[Service]
Type=oneshot
ExecStart=updategeo
; Restart=on-failure

[Install]
WantedBy=multi-user.target
